# deploy/vps/systemd/ctoa-health-live.service
[Unit]
Description=CTOA live health monitor stream
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ctoa
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/opt/ctoa/.env
ExecStart=/opt/ctoa/.venv/bin/python3 /opt/ctoa/runner/health_metrics.py --watch --interval 10 --no-publish --disk-auto-cleanup --disk-cleanup-threshold 92 --disk-cleanup-cooldown 3600
Restart=always
RestartSec=2
StandardOutput=append:/opt/ctoa/logs/health-live.log
StandardError=append:/opt/ctoa/logs/health-live.log

[Install]
WantedBy=multi-user.target
