# deploy/vps/systemd/ctoa-db.service
[Unit]
Description=CTOA PostgreSQL database (Docker)
After=docker.service network.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/ctoa
EnvironmentFile=/opt/ctoa/.env
ExecStart=/usr/bin/docker start ctoa-db
ExecStop=/usr/bin/docker stop ctoa-db
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
