# deploy/vps/systemd/ctoa-mobile-token-rotation.timer
[Unit]
Description=Rotate CTOA mobile token every 7 days

[Timer]
OnBootSec=10m
OnUnitActiveSec=7d
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
