# deploy/vps/systemd/ctoa-mythibia-news-watcher.timer
[Unit]
Description=Run CTOA Mythibia news watcher every 2 minutes

[Timer]
OnBootSec=90s
OnUnitActiveSec=2m
Persistent=true
RandomizedDelaySec=20s

[Install]
WantedBy=timers.target
