# deploy/vps/systemd/ctoa-lab-runner.service
[Unit]
Description=CTOA lab runner (autonomous project generator)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/ctoa
ExecStart=/usr/bin/python3 /opt/ctoa/runner/lab_runner.py
StandardOutput=append:/opt/ctoa/logs/lab-runner.log
StandardError=append:/opt/ctoa/logs/lab-runner.log
