# deploy/vps/systemd/ctoa-runner.service
[Unit]
Description=CTOA backlog runner tick
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/ctoa
Environment=CTOA_BACKLOG_FILE=/opt/ctoa/workflows/backlog-sprint-004.yaml
ExecStart=/usr/bin/python3 /opt/ctoa/runner/runner.py tick --agents
StandardOutput=append:/opt/ctoa/logs/runner.log
StandardError=append:/opt/ctoa/logs/runner.log
